# flask_supervisor.conf
[program:flask_nginx]
command = gunicorn app:app -b localhost:5000
directory = /app
autostart = true
autorestart = true